# build/inchiptz_1.0.1_all/etc/systemd/system/inchiptz.service
[Unit]
Description=InchiPTZ Motor Control API Server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/usr/share/inchiptz
ExecStart=/usr/bin/python3.9 /usr/share/inchiptz/api_server.py --port 192.168.25.78:502 --yaw-id 1 --pitch-id 2 --host 0.0.0.0 --port-num 50278
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# 环境变量
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=multi-user.target
